# privatebin/nginx-php.conf
[program:php-fpm]
user = user
command = sudo /usr/sbin/php-fpm7.0
autorestart = true

[program:nginx]
user = user
command = sudo /usr/sbin/nginx -g "daemon off;"
autorestart = true
